# glacier-backup.timer.example
[Unit]
Description=Run Glacier ZFS Backup weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
